# deploy/ubuntu/trading-monitor.service
[Unit]
Description=Trading Monitor (API + dashboard + supervisors)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-monitor
Environment=NODE_ENV=production
EnvironmentFile=/opt/trading-monitor/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStopSec=20
KillMode=control-group

[Install]
WantedBy=multi-user.target
